# etc/systemd/system/kevin.service
[Unit]
Description=Mount kevin's home directory.

[Service]
ExecStart=/data/mount-kevin-home.sh

[Install]
WantedBy=multi-user.target
